# packaging/aur/rewindos.install
post_install() {
  cat <<'EOF'
==> Enable the RewindOS capture daemon for your user:
      systemctl --user daemon-reload
      systemctl --user enable --now rewindos-daemon.service
==> Optional: semantic search and chat need Ollama (https://ollama.com);
    the daemon detects it automatically and pulls the embedding model.
EOF
}

post_upgrade() {
  cat <<'EOF'
==> RewindOS updated. Restart the daemon to pick up the new version:
      systemctl --user restart rewindos-daemon.service
EOF
}

pre_remove() {
  systemctl --user disable --now rewindos-daemon.service 2>/dev/null || true
}
